# Jupyter_Fabric_Testbed_Infrastructure/DNS_Infra/client_cc_files/covert_channel.service
#stored in etc/systemd/system/
[Unit]
Description=Covert Channel Python Script
After=network.target
StartLimitBurst=10
StartLimitIntervalSec=600

[Service]
ExecStart=/usr/bin/python3 /opt/covert_channel/covert_execution.py
WorkingDirectory=/opt/covert_channel
StandardOutput=journal
StandardError=journal
Restart=always
RestartSec=5s
User=root
Group=root

[Install]
WantedBy=multi-user.target
